# soar-web.service
[Unit]
Description=SOAR Web - Web Interface Service
Documentation=https://github.com/hut8/soar
After=network-online.target postgresql.service nats-server.service
Wants=network-online.target
Requires=postgresql.service
# Removed Requires=nats-server.service to prevent restart when NATS restarts
# Service will reconnect automatically via built-in retry logic

[Service]
Type=simple
User=soar
Group=soar
WorkingDirectory=/var/soar
ExecStart=/usr/local/bin/soar web --interface 127.0.0.1 --port 61225
ExecReload=/bin/kill -s HUP $MAINPID

# Environment variables
EnvironmentFile=/etc/soar/env

# Security settings
NoNewPrivileges=yes

# Restart configuration
Restart=on-failure
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Process limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target soar.target
