# soar-run.service
[Unit]
Description=SOAR Run - APRS Data Collection Service
Documentation=https://github.com/hut8/soar
After=network-online.target postgresql.service nats-server.service
Wants=network-online.target
Requires=postgresql.service nats-server.service

[Service]
Type=simple
User=soar
Group=soar
WorkingDirectory=/var/soar
ExecStart=/usr/local/bin/soar run --archive --unparsed-log /var/soar/logs/unparsed.csv
ExecReload=/bin/kill -s HUP $MAINPID

# Environment variables
EnvironmentFile=/etc/soar/env

# Security settings
NoNewPrivileges=yes
ReadWritePaths=-/var/soar

# Restart configuration
Restart=on-failure
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Process limits
LimitNOFILE=65536
MemoryMax=32G

[Install]
WantedBy=multi-user.target soar.target
