# soar-web.service
[Unit]
Description=SOAR Web - Web Interface Service
Documentation=https://github.com/your-repo/soar
After=network-online.target postgresql.service nats-server.service
Wants=network-online.target
Requires=postgresql.service nats-server.service

[Service]
Type=simple
User=soar
Group=soar
WorkingDirectory=/var/lib/soar
ExecStart=/usr/local/bin/soar web --interface 127.0.0.1 --port 61225
ExecReload=/bin/kill -s HUP $MAINPID

# Environment variables
EnvironmentFile=/etc/soar/env

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/soar/logs

# Restart configuration
Restart=on-failure
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Process limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target soar.target
