# soar-backup-base.timer
# SOAR Backup Base Timer
#
# This systemd timer automatically runs full database backups once weekly on Sunday at 2:00 AM.
# It creates a complete physical backup using pg_basebackup, compresses it, and uploads to
# cloud storage (S3-compatible).
#
# Backup details:
#   - Frequency: Weekly (Sunday 2:00 AM)
#   - Method: pg_basebackup (physical backup)
#   - Compression: gzip
#   - Storage: Cloud storage (S3/B2/etc.)
#   - Retention: 30 days (4 weekly backups)
#   - Duration: ~2-4 hours for 420GB database
#
# Combined with continuous WAL archiving, this provides point-in-time recovery capability
# for any moment within the 30-day retention period.
#
# Installation:
#   sudo cp soar-backup-base.service soar-backup-base.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable soar-backup-base.timer
#   sudo systemctl start soar-backup-base.timer
#
# Check status:
#   sudo systemctl status soar-backup-base.timer
#   sudo systemctl list-timers soar-backup-base.timer
#
# Manual execution:
#   sudo systemctl start soar-backup-base.service
#
# View logs:
#   sudo journalctl -u soar-backup-base.service -f
#   sudo tail -f /var/log/soar/backup.log

[Unit]
Description=SOAR Backup Base Timer - Weekly Full Database Backup
Requires=soar-backup-base.service

[Timer]
# Run weekly on Sunday at 2:00 AM UTC
OnCalendar=Sun *-*-* 02:00:00

# If the system was down when the timer should have triggered, run it as soon as possible
Persistent=yes

# Add some randomization to avoid exact timing conflicts (up to 30 minutes)
RandomizedDelaySec=1800

# Allow some timing drift since exact timing isn't critical for weekly backups
AccuracySec=1h

[Install]
WantedBy=timers.target
