# soar.target
# SOAR Target - Group all SOAR services together
#
# This systemd target allows you to control all SOAR services as a group while
# still maintaining individual control over each service.
#
# Usage:
#   sudo systemctl start soar      # Start all SOAR services
#   sudo systemctl stop soar       # Stop all SOAR services
#   sudo systemctl restart soar    # Restart all SOAR services
#   sudo systemctl status soar     # Check status of all SOAR services
#   sudo systemctl enable soar     # Enable all SOAR services to start on boot
#   sudo systemctl disable soar    # Disable all SOAR services from starting on boot
#
# Individual service control still works:
#   sudo systemctl start soar-web.service
#   sudo systemctl stop soar-run.service
#   sudo systemctl restart soar-pull-data.timer
#
# Installation:
#   sudo cp soar.target /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable soar.target

[Unit]
Description=SOAR - Soaring Observation And Records System
Documentation=https://github.com/your-repo/soar
After=network-online.target postgresql.service nats-server.service
Wants=soar-run.service soar-web.service soar-pull-data.timer
Requires=postgresql.service nats-server.service

[Install]
WantedBy=multi-user.target
