# soar-pull-data.timer
# SOAR Pull Data Timer
#
# This systemd timer automatically runs the `soar pull-data` command once daily at 3:00 AM.
# It downloads the latest airports, runways, and aircraft data from official sources,
# and optionally pulls devices, geocodes locations, and links home bases.
#
# Installation:
#   sudo cp soar-pull-data.service soar-pull-data.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable soar-pull-data.timer
#   sudo systemctl start soar-pull-data.timer
#
# Check status:
#   sudo systemctl status soar-pull-data.timer
#   sudo systemctl list-timers soar-pull-data.timer
#
# Manual execution:
#   sudo systemctl start soar-pull-data.service
#
# View logs:
#   sudo journalctl -u soar-pull-data.service -f

[Unit]
Description=SOAR Pull Data Timer - Run Daily Data Import
Requires=soar-pull-data.service

[Timer]
# Run daily at 3:00 AM local time
OnCalendar=*-*-* 03:00:00

# If the system was down when the timer should have triggered, run it as soon as possible
Persistent=yes

# Add some randomization to avoid all instances hitting servers at exactly the same time
RandomizedDelaySec=300

# Ensure the service isn't triggered more than once per day even if the system time changes
AccuracySec=1h

[Install]
WantedBy=timers.target
