# infrastructure/alloy.service
[Unit]
Description=Grafana Alloy
Documentation=https://grafana.com/docs/alloy/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=alloy
Group=alloy

# Working directory
WorkingDirectory=/var/lib/alloy

# Alloy command
# Configuration file is in /etc/alloy/config.alloy
ExecStart=/usr/bin/alloy run /etc/alloy/config.alloy --storage.path=/var/lib/alloy --server.http.listen-addr=127.0.0.1:12345

# Reload configuration on SIGHUP
ExecReload=/bin/kill -HUP $MAINPID

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=alloy

# Restart policy
Restart=on-failure
RestartSec=10s

# Resource limits
LimitNOFILE=65535

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/alloy

[Install]
WantedBy=multi-user.target
