# infrastructure/nats-server.service
[Unit]
Description=NATS Server with JetStream
After=network-online.target

[Service]
PrivateTmp=true
Type=simple
# Enable JetStream with file-based storage for durable message queues
# Note: Create /var/lib/nats/jetstream directory and set ownership to daemon:daemon
# mkdir -p /var/lib/nats/jetstream && chown daemon:daemon /var/lib/nats/jetstream
ExecStart=/usr/local/bin/nats-server --jetstream --store_dir=/var/lib/nats/jetstream
ExecReload=/bin/kill -s HUP $MAINPID

# The nats-server uses SIGUSR2 to trigger Lame Duck Mode (LDM) shutdown
# https://docs.nats.io/running-a-nats-service/nats_admin/lame_duck_mode
ExecStop=/bin/kill -s SIGUSR2  $MAINPID

# This should be `lame_duck_duration` + some buffer to finish the shutdown.
# By default, `lame_duck_duration` is 2 mins.
TimeoutStopSec=150

Restart=on-failure

User=daemon
Group=daemon

[Install]
WantedBy=multi-user.target
