# infrastructure/pyroscope.service
[Unit]
Description=Pyroscope Continuous Profiling Server
Documentation=https://grafana.com/docs/pyroscope/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pyroscope
Group=pyroscope

# Working directory
WorkingDirectory=/var/lib/pyroscope

# Pyroscope command
# Data is stored in /var/lib/pyroscope
# Configuration file is in /etc/pyroscope/config.yml
ExecStart=/usr/bin/pyroscope -config.file=/etc/pyroscope/config.yml

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pyroscope

# Restart policy
Restart=on-failure
RestartSec=10s

# Resource limits
LimitNOFILE=65535

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/pyroscope

[Install]
WantedBy=multi-user.target
