# infrastructure/systemd/soar-aggregate-coverage.service
[Unit]
Description=SOAR Coverage Aggregation - Aggregate position fixes into H3 coverage hexes
Documentation=https://github.com/hut8/soar
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=oneshot
User=soar
Group=soar
WorkingDirectory=/var/soar

# Aggregate coverage at all three resolutions (6, 7, 8)
# Auto-detects date range from database (aggregates from last coverage date + 1 to yesterday)
ExecStart=/usr/local/bin/soar aggregate-coverage --resolutions 6,7,8

# Environment variables
EnvironmentFile=/etc/soar/env

# Security settings
NoNewPrivileges=yes

# Process limits
LimitNOFILE=65536

# No timeout - allow aggregation to run as long as needed (especially for backlog processing)
TimeoutStartSec=infinity

[Install]
WantedBy=multi-user.target soar.target
