# infrastructure/systemd/soar-aggregate-coverage-staging.timer
# SOAR Coverage Aggregation Timer (Staging)
#
# This systemd timer automatically runs the `soar-staging aggregate-coverage` command once daily at 2:00 AM.
# It processes the previous day's position fixes into H3 hexagonal coverage zones at resolutions 6, 7, and 8
# for visualization of receiver coverage areas.
#
# Installation:
#   sudo cp soar-aggregate-coverage-staging.service soar-aggregate-coverage-staging.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable soar-aggregate-coverage-staging.timer
#   sudo systemctl start soar-aggregate-coverage-staging.timer
#
# Check status:
#   sudo systemctl status soar-aggregate-coverage-staging.timer
#   sudo systemctl list-timers soar-aggregate-coverage-staging.timer
#
# Manual execution:
#   sudo systemctl start soar-aggregate-coverage-staging.service
#
# View logs:
#   sudo journalctl -u soar-aggregate-coverage-staging.service -f

[Unit]
Description=SOAR Coverage Aggregation Timer (Staging) - Run Daily at 2 AM
Requires=soar-aggregate-coverage-staging.service

[Timer]
# Run daily at 2:00 AM local time
OnCalendar=*-*-* 02:00:00

# If the system was down when the timer should have triggered, run it as soon as possible
Persistent=yes

# Add some randomization to avoid server load spikes
RandomizedDelaySec=300

# Ensure the service isn't triggered more than once per day even if the system time changes
AccuracySec=30min

[Install]
WantedBy=timers.target soar.target
