# infrastructure/systemd/soar-aggregate-coverage-staging.service
[Unit]
Description=SOAR Coverage Aggregation (Staging) - Aggregate position fixes into H3 coverage hexes
Documentation=https://github.com/hut8/soar
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=oneshot
User=soar
Group=soar
WorkingDirectory=/var/soar

# Aggregate coverage at all three resolutions (6, 7, 8)
# Auto-detects date range from database (aggregates from last coverage date + 1 to yesterday)
ExecStart=/usr/local/bin/soar-staging aggregate-coverage --resolutions 6,7,8

# Environment variables (staging environment)
EnvironmentFile=/etc/soar/env-staging

# Security settings
NoNewPrivileges=yes

# Process limits
LimitNOFILE=65536

# Timeout (allow up to 30 minutes for coverage aggregation of large datasets)
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target soar.target
