# infrastructure/systemd/soar-backup-base.timer
# SOAR Backup Base Timer
#
# This systemd timer automatically runs full database backups every Sunday at midnight.
# It creates a complete physical backup using pg_basebackup, compresses it, and uploads to
# remote SSH server.
#
# Backup details:
#   - Frequency: Weekly (Sunday at midnight UTC)
#   - Method: pg_basebackup (physical backup)
#   - Compression: zstd
#   - Storage: Remote SSH server
#   - Retention: Keep 5 most recent backups
#   - Duration: ~2-4 hours for 420GB database
#
# Combined with continuous WAL archiving, this provides point-in-time recovery capability.
# Weekly base backups provide ~5 weeks of backup history (with keep count of 5).
#
# Installation:
#   sudo cp soar-backup-base.service soar-backup-base.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable soar-backup-base.timer
#   sudo systemctl start soar-backup-base.timer
#
# Check status:
#   sudo systemctl status soar-backup-base.timer
#   sudo systemctl list-timers soar-backup-base.timer
#
# Manual execution:
#   sudo systemctl start soar-backup-base.service
#
# View logs:
#   sudo journalctl -u soar-backup-base.service -f
#   sudo tail -f /var/log/soar/backup.log

[Unit]
Description=SOAR Backup Base Timer - Weekly Full Database Backup
Requires=soar-backup-base.service

[Timer]
# Run weekly: Every Sunday at midnight (00:00 UTC)
OnCalendar=Sun 00:00:00

# If the system was down when the timer should have triggered, run it as soon as possible
Persistent=yes

# Add some randomization to avoid exact timing conflicts (up to 30 minutes)
RandomizedDelaySec=1800

# Allow some timing drift since exact timing isn't critical for weekly backups
AccuracySec=1h

[Install]
WantedBy=timers.target
